# services.d/php-fpm/run
#!/usr/bin/with-contenv sh

su www-data

exec /usr/local/sbin/php-fpm --nodaemonize