# services.d/php-fpm/run
#!/bin/sh

su www-data

exec /usr/local/sbin/php-fpm